# heliosburn/hblog/supervisord_hblog.conf
[supervisord]
logfile = /tmp/supervisord_hblog.log

[program:hblog]
command = python -u /home/vagrant/HeliosBurn/heliosburn/hblog/hblog.py
stdout_logfile = /tmp/hblog.log
redirect_stderr = true