# install/etc/supervisor/conf.d/hbdemo.conf
[program:hbdemo]
command = python -u /opt/HeliosBurn/demo_server.py
redirect_stderr = true
