# install/etc/supervisor/conf.d/hbdjango.conf
[program:hbdjango]
command = python -u /opt/HeliosBurn/heliosburn/django/hbproject/cherrypy_launcher.py
redirect_stderr = true
