# install/etc/supervisor/conf.d/hblog.conf
[program:hblog]
command = python -u /opt/HeliosBurn/heliosburn/hblog/hblog.py
redirect_stderr = true
