# heliosburn/supervisord_conf_d/hbproxy.conf
[program:hbproxy]
directory = /opt/HeliosBurn/heliosburn
command = python -u proxy
redirect_stderr = true
